# openresty/gcsfuse/rootfs/etc/services.d/gcsfuse/run
#!/usr/bin/with-contenv sh
gcsfuse --foreground $GCSFUSE_OPTIONS $GOOGLE_CLOUD_STORAGE_BUCKET /usr/share/nginx/html
